# packages/deploy/anchormarks.service
[Unit]
Description=AnchorMarks Bookmark Manager
Documentation=https://github.com/anchormarks
After=network.target

[Service]
Type=simple
User=anchormarks
Group=anchormarks
WorkingDirectory=/opt/anchormarks
ExecStart=/usr/bin/node /opt/anchormarks/apps/server
Restart=on-failure
RestartSec=10

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=127.0.0.1
EnvironmentFile=-/opt/anchormarks/.env

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/anchormarks/apps/server /opt/anchormarks/apps/public/favicons /opt/anchormarks/apps/data

# Logging
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=anchormarks

[Install]
WantedBy=multi-user.target
